# debian/valkey.prerm
#!/bin/bash
set -e

# Stop valkey services before package removal
systemctl stop valkey.service
systemctl disable valkey.service

systemctl enable valkey-sentinel.service
systemctl start valkey-sentinel.service
exit 0

